# deepcool-cli.service
[Unit]
Description=DeepCool Digital Controller
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/deepcool-cli --mode cpu --interval 1000 --rotate 270
Restart=on-failure
RestartSec=5

# Run as root for device access
User=root

# Environment
Environment=QT_QPA_PLATFORM=offscreen

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deepcool-cli

[Install]
WantedBy=multi-user.target
